# config/supervisord-pipecat.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info

# =============================================================================
# Xvfb - Virtual Framebuffer
# Creates display :1 for headless rendering
# =============================================================================
[program:xvfb]
command=/bin/bash -c "rm -f /tmp/.X1-lock /tmp/.X11-unix/X1 2>/dev/null; exec Xvfb :1 -screen 0 %(ENV_RESOLUTION)sx24 +extension GLX +render -noreset -ac"
autostart=true
autorestart=true
priority=10
startsecs=2
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb_err.log

# =============================================================================
# Openbox - Lightweight Window Manager
# =============================================================================
[program:openbox]
command=/bin/bash -c "sleep 3 && DISPLAY=:1 openbox"
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/var/log/supervisor/openbox.log
stderr_logfile=/var/log/supervisor/openbox_err.log

# =============================================================================
# x11vnc - VNC Server for Xvfb display
# =============================================================================
[program:x11vnc]
command=/bin/bash -c "sleep 4 && x11vnc -display :1 -forever -shared -rfbport 5901 -nopw -xkb -noxrecord -noxfixes -noxdamage"
autostart=true
autorestart=true
priority=30
startsecs=6
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc_err.log

# =============================================================================
# noVNC - Web-based VNC client via websockify
# =============================================================================
[program:novnc]
command=/bin/bash -c "sleep 6 && websockify --web=/opt/noVNC %(ENV_NOVNC_PORT)s localhost:5901"
autostart=true
autorestart=true
priority=40
startsecs=8
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc_err.log

# =============================================================================
# Reachy Mini Daemon - MuJoCo Simulation
# Runs with --no-localhost-only so Pipecat bot can connect
# =============================================================================
[program:reachy-daemon]
command=/bin/bash -c "sleep 10 && python -m reachy_mini.daemon.app.main --sim --scene %(ENV_REACHY_SCENE)s --no-localhost-only"
environment=DISPLAY=":1",MUJOCO_GL="egl",__GLX_VENDOR_LIBRARY_NAME="nvidia",__NV_PRIME_RENDER_OFFLOAD="1"
autostart=true
autorestart=true
priority=50
startsecs=12
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/reachy-daemon.log
stderr_logfile=/var/log/supervisor/reachy-daemon_err.log

# =============================================================================
# NeMo Agent Toolkit Server
# Provides AI routing between Nemotron models (text, VLM, agent)
# =============================================================================
[program:nat-server]
command=/bin/bash -c "sleep 15 && cd /app/nat && nat serve --config_file src/ces_tutorial/config.yml --port %(ENV_NAT_PORT)s"
environment=PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
priority=55
startsecs=18
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/nat-server.log
stderr_logfile=/var/log/supervisor/nat-server_err.log

# =============================================================================
# Pipecat Bot - WebRTC Voice/Vision AI
# Connects to NAT for LLM routing, ElevenLabs for STT/TTS, Reachy for robot
# ICE servers (STUN/TURN) are now configured in main.py via SmallWebRTCConnection
# =============================================================================
[program:pipecat-bot]
command=/bin/bash -c "sleep 25 && cd /app/bot && python main.py --transport webrtc --host 0.0.0.0 --port 7880"
environment=DISPLAY=":1",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
priority=60
startsecs=28
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/pipecat-bot.log
stderr_logfile=/var/log/supervisor/pipecat-bot_err.log

# =============================================================================
# Jupyter Lab - Interactive Development
# =============================================================================
[program:jupyter]
command=jupyter lab --ip=0.0.0.0 --port=%(ENV_JUPYTER_PORT)s --no-browser --allow-root
directory=/workspace
environment=DISPLAY=":1"
autostart=true
autorestart=true
priority=70
startsecs=5
stdout_logfile=/var/log/supervisor/jupyter.log
stderr_logfile=/var/log/supervisor/jupyter_err.log

# =============================================================================
# Caddy - TLS Reverse Proxy for Pipecat WebRTC
# HTTPS required for browser getUserMedia (microphone/camera)
# =============================================================================
[program:caddy]
command=/bin/bash -c "sleep 28 && /usr/local/bin/caddy run --config /etc/caddy/Caddyfile"
autostart=true
autorestart=true
priority=75
startsecs=30
stdout_logfile=/var/log/supervisor/caddy.log
stderr_logfile=/var/log/supervisor/caddy_err.log

